# etc/zsh/profile.d/01_aliases
# set aliases
alias cp='cp -iv'
alias curl='curl -k'
alias grep='grep --color=auto'
alias ip='ip --color=auto'
alias ll='ls -Ahls --color=auto'
alias ls='ls -A --color=auto'
alias mkdir='mkdir -pv'
alias mv='mv -iv'
alias mvn='mvn -gs "$XDG_CONFIG_HOME"/maven/settings.xml '
alias vi='nvim'
alias rm='rm -Iv'